# .config/systemd/user/pipewire.socket.d/override.conf
[Socket]
#ExecStartPost=/bin/mount /var/lib/libvirt/qemu/pipewire-0
#ExecStopPre=/bin/umount /var/lib/libvirt/qemu/pipewire-0